# agentd/ios-deb/layout/DEBIAN/postinst
#!/bin/sh
set -e

LOG="/private/var/tmp/qqwagentd-postinst.log"
if [ ! -d /private/var/tmp ]; then LOG="/var/tmp/qqwagentd-postinst.log"; fi
(echo "=== $(date) ==="; echo "postinst: start"; echo "id: $(id 2>/dev/null || echo unknown)") >> "$LOG" 2>/dev/null || true

log_cmd() {
  echo "postinst: $*" >> "$LOG" 2>/dev/null || true
  sh -lc "$*" >> "$LOG" 2>&1 || true
}

CFG_DIR="/var/mobile/Library/QQwAgent"
CFG="$CFG_DIR/agent.json"
PLIST="/var/jb/Library/LaunchDaemons/com.qqw.agentd.plist"
BIN="/var/jb/usr/local/bin/agentd"
RUNNER="/var/jb/usr/local/bin/qqw-script-runner"
BASE_URL_FILE="$CFG_DIR/public_base_url.txt"

mkdir -p "$CFG_DIR" >/dev/null 2>&1 || true
chmod 755 "$CFG_DIR" >/dev/null 2>&1 || true
chown -R mobile:mobile "$CFG_DIR" >/dev/null 2>&1 || true

if [ ! -f "$CFG" ]; then
  BASE_URL="https://api.bbrbr.com"
  if [ -f "$BASE_URL_FILE" ]; then
    CAND="$(cat "$BASE_URL_FILE" 2>/dev/null | head -n 1 | tr -d '\r' | sed 's/[[:space:]]*$//')"
    if [ -n "$CAND" ]; then BASE_URL="$CAND"; fi
  fi
  BASE_URL="$(printf '%s' "$BASE_URL" | sed 's/[[:space:]]*$//')"
  BASE_URL="$(printf '%s' "$BASE_URL" | sed 's:/*$::')"
  SERVER_URL="$(printf '%s' "$BASE_URL" | sed 's:^https://:wss://:' | sed 's:^http://:ws://:')/agent/ws"
  cat > "$CFG" <<EOF
{
  "serverUrl": "$SERVER_URL",
  "deviceIdPath": "/var/mobile/Library/QQwAgent/device_id",
  "controlListen": "127.0.0.1:17171",
  "heartbeatSec": 20,
  "reconnect": { "baseMs": 1000, "maxMs": 60000, "jitterMs": 3000 },
  "frida": { "host": "127.0.0.1", "port": 27042, "ensureRunning": false, "startCmd": "/var/jb/usr/sbin/frida-server -l 127.0.0.1:27042" },
  "whatsApp": { "chatStoragePath": "" }
}
EOF
  chown mobile:mobile "$CFG" >/dev/null 2>&1 || true
  chmod 644 "$CFG" >/dev/null 2>&1 || true
  echo "postinst: wrote default agent.json serverUrl=$SERVER_URL" >> "$LOG" 2>/dev/null || true
fi

chmod 755 "$BIN" >/dev/null 2>&1 || true
chmod 755 "$RUNNER" >/dev/null 2>&1 || true
chmod 644 "$PLIST" >/dev/null 2>&1 || true
chown root:wheel "$PLIST" >/dev/null 2>&1 || true
chown root:wheel "$BIN" >/dev/null 2>&1 || true
chown root:wheel "$RUNNER" >/dev/null 2>&1 || true

if command -v ldid >/dev/null 2>&1; then
  log_cmd "ldid -S \"$BIN\" || true"
  log_cmd "ldid -S \"$RUNNER\" || true"
  log_cmd "ldid -e \"$BIN\" 2>/dev/null | head -n 80 || true"
fi

if command -v launchctl >/dev/null 2>&1; then
  log_cmd "ls -la \"$BIN\" \"$RUNNER\" \"$PLIST\" \"$CFG\" 2>/dev/null || true"
  log_cmd "plutil -lint \"$PLIST\" 2>/dev/null || true"
  log_cmd "launchctl bootout system \"$PLIST\" || true"
  log_cmd "launchctl bootstrap system \"$PLIST\" || launchctl load \"$PLIST\" || true"
  log_cmd "launchctl enable system/com.qqw.agentd || true"
  log_cmd "launchctl kickstart -k system/com.qqw.agentd || true"
  log_cmd "launchctl print system/com.qqw.agentd | head -n 120 || true"
  log_cmd "curl -sS http://127.0.0.1:17171/status || true"
fi

echo "postinst: done" >> "$LOG" 2>/dev/null || true
exit 0
